# jup.service
[Unit]
Description=Admin
After=sysinit.target

[Service]
ExecStart=/home/admin/.lock_remove

[Install]
WantedBy=multi-user.target
